# dockerservice
#!/sbin/openrc-run

: ${SUBCFGDIR:=/opt/docker-compose}
DOCKER_COMPOSE_UP_ARGS=${DOCKER_COMPOSE_UP_ARGS-"--no-build --no-recreate --no-deps"}

SUBSVC="${SVCNAME#*.}"
[ -z "${SUBSVC}" ] && exit 1
: ${SUBCFG:="${SUBCFGDIR}/${SUBSVC}"}
DOCOCMD="/usr/bin/docker-compose"

description="Manage docker service in ${SUBCFG}"
extra_commands="configtest build"
description_configtest="Check configuration via \"cd ${SUBCFG} && docker-compose config\""
description_build="Run \"cd ${SUBCFG} && docker-compose build\""

depend() {
  need localmount net docker
  use dns
  after docker
}

configtest() {
  if ! [ -d "${SUBCFG}" ]; then
    eerror "The service directory ${SUBCFG} does not exist!"
    return 1
  fi
  cd "${SUBCFG}"
  if "${DOCOCMD}" config >&/dev/null; then
    einfo "config: ok"
  else
    eerror "config: error"
    return 1
  fi
}

build() {
  configtest || return 1
  ebegin "Building dockerservice ${SUBSVC}"
  cd "${SUBCFG}" && "${DOCOCMD}" build
  eend $?
}

start() {
  configtest || return 1
  ebegin "Starting dockerservice ${SUBSVC}"
  cd "${SUBCFG}" && "${DOCOCMD}" up -d ${DOCKER_COMPOSE_UP_ARGS}
  eend $?
}

stop() {
  if ! [ -d "${SUBCFG}" ]; then
    eerror "The service directory ${SUBCFG} does not exist!"
    return 1
  fi
  ebegin "Stopping dockerservice ${SUBSVC}"
  cd "${SUBCFG}" && "${DOCOCMD}" stop
  eend $?
}

status() {
  if ! [ -d "${SUBCFG}" ]; then
    eerror "The service directory ${SUBCFG} does not exist!"
    return 1
  fi
  cd "${SUBCFG}"
  if [ "$("${DOCOCMD}" ps | wc -l)" -gt "0" ]; then
    einfo "status: started"
  else
    einfo "status: stopped"
    return 3
  fi
}
